# deploy/systemd/bounty-pool.service
[Unit]
Description=Bounty Pool FastAPI API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=bounty
Group=bounty
WorkingDirectory=/opt/bounty-pool/app/backend
EnvironmentFile=/opt/bounty-pool/app/backend/.env
ExecStart=/opt/bounty-pool/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
